# tests/root/etc/init/system-setup-keyboard.conf
# system-setup-keyboard - /etc/sysconfig/keyboard to xorg.conf converter
#

start on runlevel 5
stop on runlevel [!5]

console output
exec /usr/bin/system-setup-keyboard
