# tests/root/etc/init/vpnc-cleanup.conf
# vpnc-cleanup - clean up vpnc modified resolv.conf
#

start on stopped rcS

task
script
	rm -f /var/run/vpnc/pid
	RESOLV_CONF_BACKUP=/var/run/vpnc/resolv.conf-backup
	if [ -e "$RESOLV_CONF_BACKUP" ]; then
		grep '^#@VPNC_GENERATED@' /etc/resolv.conf > /dev/null 2>&1 && cat "$RESOLV_CONF_BACKUP" > /etc/resolv.conf
		rm -f -- "$RESOLV_CONF_BACKUP"
	fi
end script
